# infra/systemd/enesa-sla-monitor.service
[Unit]
Description=Enesa Automation Hub SLA Monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/enesa-automation-hub/backend
ExecStart=/usr/bin/python -m app.workers.sla_monitor
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
